# tmux/mac/launchd/com.agents-nexus.agent-memory.flush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.agent-memory.flush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.tmux/flush-events.sh</string>
    </array>

    <!-- Flush every 2 minutes -->
    <key>StartInterval</key>
    <integer>120</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/.tmux/flush-events.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.tmux/flush-events.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>AGENTS_NEXUS_DIR</key>
        <string>__AGENTS_NEXUS_DIR__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
