# tmux/mac/launchd/com.agents-nexus.arbiter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.arbiter</string>

    <!-- Pixel Dashboard bridge server (fleet state + WebSocket + permission respond).
         Reads fleet state via substrate.sh. herdr is the default backend now (arbiter
         self-defaults to it); set NEXUS_SUBSTRATE=tmux in this plist's env to pin the
         legacy tmux backend on a box that hasn't cut over. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>node</string>
        <string>__AGENTS_NEXUS_DIR__/arbiter/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__AGENTS_NEXUS_DIR__/arbiter</string>

    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-arbiter.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-arbiter.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>AGENTS_NEXUS_DIR</key>
        <string>__AGENTS_NEXUS_DIR__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string>
        <key>TMUX_SESSION</key>
        <string>agents</string>
    </dict>
</dict>
</plist>
