# tmux/linux/systemd/arbiter.service
[Unit]
Description=Arbiter — tmux-to-dashboard WebSocket bridge
After=tmux-agents-session.service
Requires=tmux-agents-session.service

[Service]
Type=simple
WorkingDirectory=__AGENTS_NEXUS_DIR__/arbiter
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
Environment=HOME=__HOME__
Environment=PORT=8420
Environment=TMUX_SESSION=agents

[Install]
WantedBy=default.target
