# tmux/linux/systemd/crash-breadcrumb.service
[Unit]
Description=Crash breadcrumb — fsync'd temp/power/load snapshot to survive spontaneous power-off

[Service]
Type=simple
ExecStart=__HOME__/.tmux/crash-breadcrumb.sh
Restart=always
RestartSec=5
Nice=15
IOSchedulingClass=idle

[Install]
WantedBy=default.target
