# tmux/linux/systemd/nightly-spark.service
[Unit]
Description=Nightly Spark reclaim + decision synthesis
After=nightly-repo-sync.service

[Service]
Type=oneshot
ExecStart=/bin/bash __AGENTS_NEXUS_DIR__/spark/scripts/spark-pipeline.sh
# After the index re-syncs, refresh the spawnable-repo description cache the
# Slack orchestrator merges in (scripts/spark-summary.py -> ~/.tmux/spark-summaries.json).
# Leading `-` = best-effort: a failure here never fails the nightly Spark run.
ExecStartPost=-__AGENTS_NEXUS_DIR__/spark/.venv/bin/python __AGENTS_NEXUS_DIR__/scripts/spark-summary.py
Environment=HOME=__HOME__
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
