# tmux/linux/systemd/nightly-vault-commit.service
[Unit]
Description=Auto-commit and push vault changes
After=nightly-repo-sync.service

[Service]
Type=oneshot
ExecStart=/bin/bash __AGENTS_NEXUS_DIR__/scripts/vault-commit.sh
Environment=HOME=__HOME__
Environment=VAULT_DIR=__HOME__/vault
StandardOutput=journal
StandardError=journal
