# tmux/linux/systemd/slack-bridge.service
[Unit]
Description=Slack bridge — two-way #nexus <-> agent control (Socket Mode + /notify)
# Soft ordering: the bridge delivers keystrokes into the agents tmux session via
# ~/.tmux/agent-send.sh. Wants (not Requires) so the bridge still starts and
# no-ops deliveries if the session is down.
After=tmux-agents-session.service
Wants=tmux-agents-session.service

[Service]
Type=simple
WorkingDirectory=__AGENTS_NEXUS_DIR__/slack-bridge
# Secrets come through the resolver chain (scripts/secrets/secret-run.sh), which self-defaults
# to env,doppler with the historical nexus/prd scope — so this keeps pulling SLACK_* from
# Doppler exactly as before, but a box on a different backend (AWS SM, vault) works too.
# secret-run resolves the named vars, exports them, then execs node (preserving PID for
# Type=simple); index.js's .env loader stays as a harmless fallback for anything not declared.
# DOPPLER_BIN is passed absolute (below) so the doppler backend is found under systemd's PATH.
# NATS_URL/USER/PASS join the pull for the NATS A2A transport (auth-callout user/pass +
# broker addr). Unused vars resolve empty and stay UNSET (secret-run never exports ""), so a
# slack-transport box is unaffected; index.js only reads them when NEXUS_BUS_TRANSPORT=nats.
ExecStart=__AGENTS_NEXUS_DIR__/scripts/secrets/secret-run.sh --project nexus --config prd SLACK_BOT_TOKEN SLACK_APP_TOKEN SLACK_NEXUS_CHANNEL SLACK_AGENTS_CHANNEL SLACK_BUS_ENABLED NATS_URL NATS_USER NATS_PASS -- __NODE_BIN__ index.js
# Restart only on a CRASH (non-zero exit). The bridge boot-guards to a clean
# exit 0 when SLACK_BOT_TOKEN / SLACK_APP_TOKEN are unset; on-failure (not
# always) keeps an un-provisioned box from thrashing — mirrors the launchd
# KeepAlive{SuccessfulExit:false} on the Mac plist.
Restart=on-failure
RestartSec=10
Environment=HOME=__HOME__
Environment=AGENTS_NEXUS_DIR=__AGENTS_NEXUS_DIR__
# Absolute Doppler path for the doppler backend (systemd's PATH won't find a shimmed CLI).
Environment=DOPPLER_BIN=__DOPPLER_BIN__
# tmux must be on PATH so the bridge's child agent-send.sh can deliver keystrokes
# (node + secret-run are invoked by absolute path in ExecStart, so they need no PATH entry).
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
