# tmux/linux/systemd/optional/overseer-reap.service
[Unit]
Description=Overseer reaper — checkpoint + close agents idle >4h (opt-in)
# OPT-IN: lives under systemd/optional/ so tmux/linux/install.sh does NOT
# auto-enable it (it would otherwise start closing windows). Enable deliberately
# with the steps in docs/overseer.md.

[Service]
Type=oneshot
WorkingDirectory=__AGENTS_NEXUS_DIR__
ExecStart=__AGENTS_NEXUS_DIR__/scripts/overseer-reap.sh
Environment=HOME=__HOME__
Environment=AGENTS_NEXUS_DIR=__AGENTS_NEXUS_DIR__
# This box is unattended for days — prune everything idle, command post included
# (checkpointed first). Drop a name/%pane into ~/.tmux/overseer-exclude to spare it.
Environment=REAP_ALL=1
# tmux + claude must be reachable for the reaper and its checkpoint step.
# `claude` is installed via fnm at ~/.local/share/fnm/aliases/default/bin (the
# stable `default` alias symlink), NOT ~/.local/bin — so it MUST be on PATH here
# or checkpoint-transcript.sh fails its `command -v claude` guard, exits 0 (fail
# open), and agents get reaped WITHOUT a checkpoint. (Same fnm-PATH gotcha as the
# display-popup window-launch bug.)
Environment=PATH=__HOME__/.local/share/fnm/aliases/default/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
