# launchd/optional/com.agents-nexus.overseer-reap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.overseer-reap</string>

    <!-- OPT-IN job: closes agents idle >4h (checkpoint first). Lives under
         launchd/optional/ so `launchd:install:all` does NOT sweep it in — you
         enable it deliberately with `task launchd:install:overseer-reap`.
         The script itself is fail-open and idempotent.

         PATH: a login shell (`zsh -l`) sources ~/.zprofile (Homebrew tmux) but
         NOT ~/.zshrc — and nvm (which provides `claude`) initializes in ~/.zshrc.
         Without `claude`, checkpoint-transcript.sh fail-opens and agents get
         reaped WITHOUT a checkpoint. So source nvm explicitly here to activate
         the default node + claude, version-independently. (Mirrors the Linux
         fnm-PATH fix on overseer-reap.service.) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>export NVM_DIR="$HOME/.nvm"; [ -s "$NVM_DIR/nvm.sh" ] &amp;&amp; . "$NVM_DIR/nvm.sh"; AGENTS_NEXUS_DIR=__AGENTS_NEXUS_DIR__ __AGENTS_NEXUS_DIR__/scripts/overseer-reap.sh</string>
    </array>

    <!-- Every 15 minutes. -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-overseer-reap.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-overseer-reap.log</string>
</dict>
</plist>
