# launchd/com.agents-nexus.guilty-spark.nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.guilty-spark.nightly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__AGENTS_NEXUS_DIR__/spark/scripts/spark-pipeline.sh</string>
    </array>

    <!-- Nightly at 22:00 local (Mac typically awake; avoids on-wake CPU spike) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>22</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>__AGENTS_NEXUS_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/alex.persinger@getgarner.com</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-guilty-spark-nightly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-guilty-spark-nightly.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
