# launchd/com.agents-nexus.obs-tidy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.obs-tidy</string>

    <key>ProgramArguments</key>
    <array>
        <string>__AGENTS_NEXUS_DIR__/scripts/obs-tidy</string>
    </array>

    <!-- Daily at 06:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>6</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/alex.persinger@getgarner.com</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/alex.persinger@getgarner.com/.local/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-obs-tidy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-obs-tidy.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
