# launchd/com.agents-nexus.obs-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.obs-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>__AGENTS_NEXUS_DIR__/scripts/obs-digest</string>
    </array>

    <!-- Weekdays at 07:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/alex.persinger@getgarner.com</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/alex.persinger@getgarner.com/.local/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-obs-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-obs-digest.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
