# launchd/com.agents-nexus.gl-reviews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents-nexus.gl-reviews</string>

    <!-- Run as a login shell so ~/.zprofile is sourced (GITLAB_TOKEN, PATH, etc.)
         Guard: skip if waking after noon (laptop was sleeping at scheduled time) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>cd __AGENTS_NEXUS_DIR__ &amp;&amp; [[ $(date +%H) -lt 12 ]] &amp;&amp; task gl:reviews:daily || echo "Skipped — past noon on wake"</string>
    </array>

    <!-- Weekdays at 08:05 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/agents-nexus-gl-reviews.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agents-nexus-gl-reviews.log</string>
</dict>
</plist>
