# launchd/personal/com.alexpersinger.svc-chatbot-rebase.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexpersinger.svc-chatbot-rebase</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/alex.persinger@getgarner.com/bin/svc-chatbot-rebase.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>13</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>13</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>13</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>13</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>13</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/alex.persinger@getgarner.com/Library/Logs/svc-chatbot-rebase.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/alex.persinger@getgarner.com/Library/Logs/svc-chatbot-rebase.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/alex.persinger@getgarner.com</string>
    </dict>
</dict>
</plist>
